# com.ts-question.telegram.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ts-question.telegram</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jane-m4/Desktop/claude_dev/TypeScript 专项题库/send_ts_to_telegram.py</string>
    </array>

    <!-- 每天 10:05 触发（Claude 定时任务 10:00 生成文件，留 5 分钟缓冲） -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <!-- 日志路径 -->
    <key>StandardOutPath</key>
    <string>/Users/jane-m4/Desktop/claude_dev/TypeScript 专项题库/logs/telegram.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jane-m4/Desktop/claude_dev/TypeScript 专项题库/logs/telegram.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
