# com.ts-question.bot-responder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ts-question.bot-responder</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jane-m4/Desktop/claude_dev/TypeScript 专项题库/ts_bot_responder.py</string>
    </array>

    <!-- 每 60 秒轮询一次 -->
    <key>StartInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/jane-m4/Desktop/claude_dev/TypeScript 专项题库/logs/bot-responder.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jane-m4/Desktop/claude_dev/TypeScript 专项题库/logs/bot-responder.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
